# salt/roots/salt/supervisor/coveragereport.conf
[program:coveragereport]
command=python -m SimpleHTTPServer 8001
stopsignal=INT
user=vagrant
directory=/vagrant/coveragereport
autostart=false
autorestart=false
